# infra/server/systemd/remlab-watchdog.timer
[Unit]
Description=Daily remlab disk watchdog

[Timer]
OnCalendar=*-*-* 07:00
Persistent=true

[Install]
WantedBy=timers.target
