# infra/server/systemd/remlab-db-backup.timer
[Unit]
Description=Nightly remlab DB backup

[Timer]
OnCalendar=*-*-* 03:30
Persistent=true

[Install]
WantedBy=timers.target
