# infra/server/systemd/remlab-watchdog.service
[Unit]
Description=remlab disk usage watchdog

[Service]
Type=oneshot
ExecStart=/opt/remlab/scripts/disk-watchdog.sh 80
